# engineering_site.service
[Unit]
Description=Engineering Site Gunicorn Daemon
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=notify
User=django
Group=www-data
RuntimeDirectory=gunicorn
WorkingDirectory=/home/django/projects/engineering_site
EnvironmentFile=/home/django/projects/engineering_site/.env
Environment="DJANGO_SETTINGS_MODULE=engineering_site.settings_production"
Environment="PYTHONUNBUFFERED=1"

ExecStart=/home/django/projects/engineering_site/venv/bin/gunicorn \
    --config /home/django/projects/engineering_site/gunicorn_config.py \
    engineering_site.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5s

# Ограничения безопасности
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/django/projects/engineering_site/media
ReadWritePaths=/home/django/projects/engineering_site/logs

[Install]
WantedBy=multi-user.target

